# osmind-java/agent-macos/launchd/dev.osmind.escollector.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.osmind.escollector</string>
    <key>ProgramArguments</key>
    <array>
        <string>__OSMIND_ROOT__/agent-macos/build/osmind-es-collector</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OSMIND_STORE</key>
        <string>__OSMIND_STORE__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/osmind-escollector.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/osmind-escollector.err.log</string>
</dict>
</plist>
